# systemd/.config/systemd/user/taskw-notify.timer
[Unit]
Description=Taskwarrior check for upcoming tasks

[Timer]
OnCalendar=*:0/3
Unit=taskw-notify.service

[Install]
WantedBy=default.target
